# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:simple-redirect]
command=/var/simple-redirect/env/bin/gunicorn --forwarded-allow-ips="*" -b 0.0.0.0:8000 server:app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
